# ksm-throttler.service.in
[Unit]
Description=KSM throttling daemon
Documentation=https://@PACKAGE_URL@

[Service]
ExecStart=@libexecdir@/@PACKAGE_NAME@/@TARGET@ -log debug
Restart=always

[Install]
WantedBy=multi-user.target
