# vc-throttler.service.in
[Unit]
Description=Virtcontainers based KSM throttling
Documentation=https://@PACKAGE_URL@
Requires=ksm-throttler.service

[Service]
ExecStart=@libexecdir@/@PACKAGE_NAME@/trigger/virtcontainers/vc -log debug
Restart=always

[Install]
WantedBy=multi-user.target
